# supervisor/multi-ticc.ini
[program:multi-ticc]
command=/opt/multi-ticc/bootstrap.sh /opt/multi-ticc /etc/multi-ticc.yml
autostart=true
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
user=tsc
